# .config/systemd/user/mako.service
[Unit]
Description=Notification Daemon

[Service]
Type=simple
Restart=always
ExecStart=/home/j/bin/mako
