# .config/systemd/user/waybar.service
[Unit]
Description=A status bar

[Service]
Type=simple
Restart=always
ExecStart=/home/j/bin/waybar
