# .config/systemd/user/xss-lock.service
[Unit]
Description=Screen Locker

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/xss-lock slock
