# .config/sv/waybar/finish
#!/bin/sh

exec killall -u $(id -un) waybar
