# .config/sv/waybar/run
#!/bin/sh

. ${XDG_RUNTIME_DIR:-/run/user/$(id -u)}/sway_services.env
LANG=C exec setsid /usr/bin/waybar
