# .config/sv/mako/finish
#!/bin/sh

exec killall -u $(id -un) mako
