# .config/sv/swaybg/run
#!/bin/sh

. ${XDG_RUNTIME_DIR:-/run/user/$(id -u)}/sway_services.env
exec setsid swaybg -o '*' -i ~/.wallpaper
