# etc/udev/rules.d/50-msergosculpt.rules
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="07a5", RUN+="/usr/bin/sudo -u j /home/j/bin/consistent-keyboards.sh"
